# monthly-backup.service
[Unit]
Description=Run monthly backup script

[Service]
Type=simple
ExecStart=/bin/bash /usr/bin/monthly_backup
