# system-backup-monthly.service
[Unit]
Description=Run monthly system backup script

[Service]
Type=simple
ExecStart=/bin/bash /usr/bin/system_backup_monthly
