# system-backup.service
[Unit]
Description=Run monthly system backup script

[Service]
Type=simple
ExecStart=/bin/bash /home/james/SYNC/bin/system_backup
