# system-backup-weekly.service
[Unit]
Description=Run weekly system backup script

[Service]
Type=simple
ExecStart=/bin/bash /usr/bin/system_backup_weekly
