# weekly-backup.service
[Unit]
Description=Run weekly backup script

[Service]
Type=simple
ExecStart=/bin/bash /usr/bin/weekly_backup
